# roles/certbot/templates/certbot-renewal.timer.j2
[Unit]
Description=Certbot Renewal Timer

[Timer]
OnCalendar=*-*-* 03:30:00
RandomizedDelaySec=1h
Persistent=true

[Install]
WantedBy=timers.target
